# features/attest/system_files/usr/lib/systemd/system/coolbx-device-secret.service
[Unit]
Description=Coolbx OS — genereer per-toestel HMAC-secret (first-boot)
Documentation=https://coolbx.be/
ConditionPathExists=!/etc/coolbx/device-secret
Before=coolbx-attestd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/coolbx-gen-device-secret
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
