# features/attest/system_files/usr/lib/systemd/system/coolbx-device-secret.service
[Unit]
Description=Coolbx OS — per-toestel HMAC-secret genereren + TPM2-sealen (B3.e)
Documentation=https://coolbx.be/
# Draait zolang er geen sealed credential is: first-boot (genereren) én het
# upgrade-pad (bestaand plain secret in-place sealen zodra een TPM bruikbaar
# is). Zonder TPM is de run een goedkope no-op die de toestelklasse bijschrijft.
ConditionPathExists=!/etc/coolbx/device-secret.cred
Before=coolbx-attestd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/coolbx-gen-device-secret
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
