# system_files/usr/lib/systemd/system/coolbx-firstboot-user.service
[Unit]
Description=Coolbx OS first-boot dev test user
ConditionPathExists=/etc/coolbx/enable-firstboot-user
ConditionPathExists=!/var/lib/coolbx/firstboot-user-done
After=local-fs.target
Before=gdm.service systemd-user-sessions.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/coolbx-firstboot-user.sh

[Install]
WantedBy=multi-user.target